# openwrt/files/etc/hotplug.d/input/30-mjpg-streamer
case "$ACTION" in
	add)
		# start mjpg-streamer
		/etc/init.d/mjpg-streamer start
		;;
	remove)
		 # stop mjpg-streamer
		/etc/init.d/mjpg-streamer stop
		;;
esac

